# 02_Terraform/playbooks/dashboard/semaphore.service
# systemd unit for a native (non-Docker) Semaphore install.
# Install to /etc/systemd/system/semaphore.service, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now semaphore
[Unit]
Description=Ansible Semaphore (scheduler + run-history dashboard)
Documentation=https://docs.semaphoreui.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=semaphore
Group=semaphore
ExecStart=/usr/local/bin/semaphore server --config /etc/semaphore/config.json
Restart=on-failure
RestartSec=5
# Semaphore shells out to ansible/git; give it a sane PATH.
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# The backup playbook reads the Proxmox token from the environment.
# Keep the real secret in /etc/semaphore/semaphore.env (chmod 600), NOT in git.
EnvironmentFile=-/etc/semaphore/semaphore.env

[Install]
WantedBy=multi-user.target
